# supervisor/crontab
*/5 * * * * supervisorctl start bloomchecker